# systemd/jobd.service
[Unit]
Description=jobd queue daemon
After=mariadb.service

[Service]
User=nobody
Group=nobody
PermissionsStartOnly=true
Restart=on-failure
ProtectSystem=full
PrivateDevices=true
ExecStartPre=/usr/bin/install -d -o nobody /var/log/jobd
ExecStart=/usr/bin/jobd --config /etc/jobd.conf

[Install]
WantedBy=multi-user.target